# packages/ambient-edge-server/ambient_edge_server-1.8.0-py3-none-any.whl/ambient_edge_server/templates/ambient_edge_server.service.jinja2
[Unit]
Description={{ service.description }}
After=network.target

[Service]
Type=simple
User={{ service.user }}
Group={{ service.group }}
ExecStart={{ service.exec_start }}
Restart=on-failure
{% if service.environment %}
{% for key, value in service.environment.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
{% endif %}
# WorkingDirectory={{ service.working_directory }}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target